# com.youruser.pr-review-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.youruser.pr-review-dashboard</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/youruser/projects/pr-review-dashboard/pr-review-dashboard</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/youruser/projects/pr-review-dashboard</string>
    <key>KeepAlive</key><true/>
    <key>RunAtLoad</key><true/>
    <key>StandardOutPath</key>
    <string>/Users/youruser/projects/pr-review-dashboard/dashboard.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/youruser/projects/pr-review-dashboard/dashboard.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>DB_PATH</key>
        <string>/Users/youruser/projects/pr-review-dashboard/leaderboard.db</string>
    </dict>
</dict>
</plist>
